# clockwork-orange.service
[Unit]
Description=Clockwork Orange - Wallpaper Service
After=graphical-session.target

[Service]
Type=simple
Environment=DISPLAY=:0
# Uses system Python 3 - ensure PyYAML is installed system-wide
# For conda/virtual environments, update the path to your Python interpreter
ExecStart=/usr/bin/python3 /home/nverenin/git/clockwork-orange/clockwork-orange.py
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
